# templates/halko-dbusunit.service
[Unit]
Description=Halko VPN and Power Management Service
After=network-online.target dbus.service
Wants=network-online.target
Requires=dbus.service
StartLimitIntervalSec=0

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/halko/dbusunit -c /etc/opt/halko.cfg
Restart=on-failure
RestartSec=5s
TimeoutStopSec=10s
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
